# docs/busi_route.ini
[program:busi_route]
command=/Volumes/data/devProj/virtual_proj/osspushenv/bin/celery -A busi_route worker -n busiroute --loglevel=info
redirect_stderr=true
stdout_logfile=/Volumes/data/tmp/msg_busi_route.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
directory=/Volumes/data/devProj/OSS_Push_Svr
user=jackyzhou
