# docs/celery_flower.ini
[program:celery_flower]
command=/Volumes/data/devProj/virtual_proj/dev1/bin/celery  flower --broker=amqp://guest:guest@localhost//
redirect_stderr=true
stdout_logfile=/Volumes/data/tmp/celery_flower.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
directory=/Volumes/data/devProj/virtual_proj/proj
user=jackyzhou
