# docs/celery_svr.ini
[program:celery_svr]
command=/Volumes/data/devProj/virtual_proj/dev1/bin/celery -A proj worker  --loglevel=info 
redirect_stderr=true
stdout_logfile=/Volumes/data/tmp/celery_svr.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
directory=/Volumes/data/devProj/virtual_proj
user=jackyzhou
