# docs/msg_push.ini
[program:msg_push]
command=/Volumes/data/devProj/virtual_proj/osspushenv/bin/celery -A msg_push worker -n msgpush --loglevel=info
redirect_stderr=true
stdout_logfile=/Volumes/data/tmp/msg_push.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
directory=/Volumes/data/devProj/OSS_Push_Svr
user=jackyzhou
