# docs/push_svr_api.ini
[program:push_svr_api]
command=/Volumes/data/devProj/virtual_proj/osspushenv/bin/gunicorn -c /Volumes/data/devProj/OSS_Push_Svr/gunicorn.py wsgi:app
redirect_stderr=true
stdout_logfile=/tmp/push_svr_api.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=10
directory=/Volumes/data/devProj/OSS_Push_Svr
user=jackyzhou
